# plaud-digest.timer
[Unit]
Description=Run the Plaud weekly digest every Sunday evening

[Timer]
OnCalendar=Sun 17:00
Persistent=true

[Install]
WantedBy=timers.target
